# etc/rc.d/rc.local
setconf LIBPATH $(getconf LIBPATH):/usr/pkg/lib
/sbin/devc-serusb -d vid=0x0403,did=0x6010,module=ftdi
